# .docker/services/app/etc/supervisor.d/base.ini
[inet_http_server]
port=*:9009
username=app
password=app

[supervisord]
pidfile=/tmp/supervisord.pid
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0

[program:nginx]
command=nginx -g 'pid /tmp/nginx.pid;daemon off;'
autostart=true
autorestart=true
priority=32
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stdout_events_enabled=true
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stderr_events_enabled=true

[program:php-fpm]
command=php-fpm --nodaemonize --allow-to-run-as-root
autostart=true
autorestart=true
priority=16
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stdout_events_enabled=true
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stderr_events_enabled=true
