# .docker/services/app/etc/supervisor.d/develop.ini
[program:node]
command=/bin/sh -c 'yarn install && yarn encore dev-server --host 0.0.0.0 --port 8080'
autostart=true
autorestart=true
priority=8
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stdout_events_enabled=true
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stderr_events_enabled=true

[program:php-fpm]
command=/bin/sh -c 'composer install && php-fpm --nodaemonize --allow-to-run-as-root'
